# cinexis_remote/rootfs/etc/services.d/cinexis/run
#!/usr/bin/with-contenv sh
set -eu
exec /usr/bin/cinexis-entrypoint.sh
