# cinexis_remote/rootfs/etc/services.d/cinexis/run
#!/usr/bin/with-contenv bash
exec /usr/bin/cinexis-entrypoint.sh
